# systemd/vxlan-watchdog.service
[Unit]
Description=VXLAN controller sync watchdog (alerts on persistent controller divergence)
Documentation=https://github.com/PoemaIX/vxlan-controller
After=vxlan-controller-client.service
Wants=vxlan-controller-client.service

[Service]
Type=simple
ExecStart=/usr/local/bin/vxlan-watchdog.sh
Restart=always
RestartSec=10
# Tuning / alerting (edit as needed):
# Environment=WATCHDOG_INTERVAL=30
# Environment=WATCHDOG_THRESHOLD=3
# Environment=WATCHDOG_ALERT_CMD=curl -sS -X POST -d @- https://hooks.example/vxlan

[Install]
WantedBy=multi-user.target
